# packages/docker/twenty-app-dev/rootfs/etc/s6-overlay/s6-rc.d/twenty-worker/run
#!/bin/sh
cd /app/packages/server
exec yarn worker:prod
